# config/airpi-ap3000e/luci-app-Airpifanctrl/root/etc/init.d/Airpifanctrl
#!/bin/sh /etc/rc.common
START=99

printMsg() {
    local msg="$1"
    logger -t "${PROGRAM}" "${msg}"
} 

restart(){                                                                                                       
  /usr/bin/fancts.sh &
  printMsg "Airpi风扇控制模块开始工作！"
  cp /usr/bin/footer.htm /usr/lib/lua/luci/view/themes/argon/footer.htm
  cp /usr/bin/wps /etc/rc.button/wps
}

start(){
  /usr/bin/fancts.sh &
  printMsg "Airpi风扇控制模块开始工作！"
  cp /usr/bin/footer.htm /usr/lib/lua/luci/view/themes/argon/footer.htm
  cp /usr/bin/wps /etc/rc.button/wps
}
